# ethernet/crontab
*/2 * * * * /home/pi/local/bin/monitor_eth.sh
